# etc/crontab
MAILTO=uid-1@example.com

*/5	*	*	*	*	/var/www/vhosts/fra-flugplan.de/httpdocs/fra-flugplan/getflights.sh
