# etc/crontab
MAILTO=flederwiesel@fra-flugplan.de

*/5	*	*	*	*	/var/www/vhosts/fra-flugplan.de/httpdocs/fra-flugplan/getflights.sh
