# etc/crontab
MAILTO=hausmeister@flederwiesel.com

# M H d m u
#0 0 * * *	cron.sh
*/5 	* 	* 	* 	* 	php -f /home/flederwi/public_html/fra-schedule/getflights.php >> /home/flederwi/logs/getflights.log
